# src/crontab
1 * * * * $HOME/p/alarmclock/src/updater.sh
