# src/crontab
*/5 * * * * $HOME/p/alarmclock/src/updater.sh
